# deploy/systemd/kairos-deadman.service
[Unit]
Description=KAIROS dead-man's-switch (alerts #kairos-deadman)

[Service]
Type=oneshot
User=kairos
Group=kairos
EnvironmentFile=/etc/kairos/kairos.env
ExecStart=/opt/kairos/deploy/deadman.sh
