# deploy/systemd/kairos-digest.timer
[Unit]
Description=Run KAIROS digest every 2 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=2h
AccuracySec=1min

[Install]
WantedBy=timers.target
